# templates/mount-data-disks.service
[Unit]
Description=Mount data disks

[Service]
Type=oneshot
ExecStart={{mdd_mount_path}}

[Install]
WantedBy=multi-user.target
