# files/mount-data-disks.service
[Unit]
Description=Mount data disks

[Service]
Type=oneshot
ExecStart=/usr/bin/mount-data-disks

[Install]
WantedBy=multi-user.target
